# config/non-docker/services/launchd/com.cryptobot.trade.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cryptobot.trade</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>-m</string>
        <string>uvicorn</string>
        <string>trade.main:app</string>
        <string>--host</string>
        <string>0.0.0.0</string>
        <string>--port</string>
        <string>8004</string>
        <string>--workers</string>
        <string>2</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>WorkingDirectory</key>
    <string>/opt/cryptobot</string>
    <key>StandardErrorPath</key>
    <string>/opt/cryptobot/logs/trade_service_error.log</string>
    <key>StandardOutPath</key>
    <string>/opt/cryptobot/logs/trade_service.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/opt/cryptobot</string>
    </dict>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>ProcessType</key>
    <string>Background</string>
    <key>Dependency</key>
    <array>
        <dict>
            <key>Label</key>
            <string>com.cryptobot.auth</string>
            <key>Required</key>
            <true/>
        </dict>
        <dict>
            <key>Label</key>
            <string>com.cryptobot.strategy</string>
            <key>Required</key>
            <true/>
        </dict>
    </array>
</dict>
</plist>