# config/non-docker/services/systemd/cryptobot-data.service
[Unit]
Description=CryptoBot Data Service
After=network.target redis.service cryptobot-auth.service
Wants=network.target redis.service
Requires=cryptobot-auth.service

[Service]
Type=simple
User=cryptobot
Group=cryptobot
WorkingDirectory=/opt/cryptobot
ExecStart=/usr/bin/python3 -m uvicorn data.main:app --host 0.0.0.0 --port 8001 --workers 2
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cryptobot-data
Environment="PYTHONPATH=/opt/cryptobot"
EnvironmentFile=/opt/cryptobot/.env

# Security settings
PrivateTmp=true
ProtectSystem=full
NoNewPrivileges=true
ProtectHome=true

[Install]
WantedBy=multi-user.target